# deploy/arch-repo/installation/files/etc/systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin cantoro --noclear %I $TERM